# scripts/supervisor.conf
[supervisord]
nodaemon=true

[program:nginx]
command=/bin/bash -c "/usr/sbin/nginx -g 'daemon off;'"

[program:u235core]
directory=/u235core
command=/bin/bash -c "/usr/bin/node /u235core/main.js"
